# root/etc/s6-overlay/s6-rc.d/init-cups-config/run
#!/usr/bin/with-contenv sh

rc-service cups start
rc-service cupsd start