# root/etc/s6-overlay/s6-rc.d/init-cups-config/run
# run文件内容
#!/command/with-contenv bash
lpadmin -p AirPrint -v socket://printer.local:9100 -E
lpstat -v  # 验证打印机
