# root/etc/s6-overlay/s6-rc.d/init-cups-config/run
# run 文件内容
#!/command/with-contenv bash
# 示例：初始化CUPS打印配置
lpadmin -p AirPrint -v socket://printer.local:9100 -E
echo "CUPS配置完成"
