# root/etc/s6-overlay/s6-rc.d/init-cups-config/run
#!/usr/bin/with-contenv sh

# 创建必要目录
mkdir -p \
    /var/run/cups \
    /etc/cups/ssl

# 设置CUPS配置文件路径
if [ ! -f /etc/cups/cupsd.conf ]; then
    cp -rn /etc/cups/cupsd.conf.default /etc/cups/cupsd.conf
fi

# 修复权限
chown -R root:lpadmin /etc/cups
chmod 0755 /var/run/cups

# 启动CUPS服务
rc-service cupsd start 2>/dev/null
rc-service cupsd restart 2>/dev/null

# 健康检查等待启动
while ! nc -z localhost 631; do
    sleep 1
done

# 持久化打印机配置
if [ -d /config/printers ]; then
    cp -rf /config/printers/* /etc/cups/
    lpadmin -E -p $(ls /etc/cups/ppd/*.ppd | sed 's/.*\///;s/.ppd//')
fi

# CUPS 权限处理核心逻辑
echo "**** 配置CUPS权限 ****" && \
mkdir -p \
    /var/run/cups \
    /etc/cups/ssl \
    /config/printers/ppd && \

# 权限继承处理
if [ -d /config/printers ]; then
    rsync -rtD --exclude '*.tmp' /config/printers/ /etc/cups/
fi && \

# 关键目录权限设置
chown -R root:lpadmin /etc/cups && \
chmod 0750 /var/run/cups && \
chmod 0700 /etc/cups/ssl && \
chmod 0644 /etc/cups/cupsd.conf* && \

# PPD文件特殊权限处理
find /etc/cups/ppd -name "*.ppd" -exec chmod 0664 {} \; && \
find /etc/cups/ppd -name "*.ppd" -exec chown root:lpadmin {} \; && \

# 运行时目录权限修复
( [ -d /var/spool/cups ] || mkdir -p /var/spool/cups ) && \
chown root:lp /var/spool/cups && \
chmod 1770 /var/spool/cups && \

# 容器用户权限兼容处理
if id -u abc >/dev/null 2>&1; then
    usermod -aG lpadmin abc && \
    find /config/printers -exec chown abc:abc {} \;
fi